# config/crashmanager.service.in
###############################################################################
#
#  Coredump manager
#
###############################################################################
[Unit]
Description=Coredump manager service
Before=dlt-system.service
Wants=dlt-system.service

[Service]
Type=simple
ExecStart=@install_prefix@/bin/crashmanager
WatchdogSec=30

[Install]
WantedBy=multi-user.target
